# ops/systemd/claimcoin-flaresolverr.service
[Unit]
Description=ClaimCoin FlareSolverr Helper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/state/flaresolverr-exp/src
Environment=HOME=/root
ExecStartPre=/bin/mkdir -p /root/.openclaw/workspace/state/flaresolverr-exp/logs
ExecStart=/root/.openclaw/workspace/state/flaresolverr-exp/src/.venv/bin/python src/flaresolverr.py
Restart=always
RestartSec=2
StandardOutput=append:/root/.openclaw/workspace/state/flaresolverr-exp/logs/run.log
StandardError=append:/root/.openclaw/workspace/state/flaresolverr-exp/logs/run.log

[Install]
WantedBy=multi-user.target
