# ops/systemd/claimcoin-antibot.service
[Unit]
Description=ClaimCoin AntiBot Solver API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/projects/indra-api-hub
Environment=HOME=/root
Environment=PYTHONPATH=/root/.openclaw/workspace/projects/indra-api-hub/integrations/antibot-image-solver/src
Environment=ANTIBOT_OCR_PROFILE=fast
ExecStart=/root/.openclaw/workspace/projects/indra-api-hub/.venv/bin/python -m antibot_image_solver serve --host 127.0.0.1 --port 8010
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
