# ops/systemd/claimcoin-runloop.service
[Unit]
Description=ClaimCoin Autoclaim Run Loop
After=network-online.target claimcoin-antibot.service claimcoin-flaresolverr.service
Wants=network-online.target claimcoin-antibot.service claimcoin-flaresolverr.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/projects/claimcoin-autoclaim
Environment=HOME=/root
Environment=PYTHONPATH=/root/.openclaw/workspace/projects/claimcoin-autoclaim/src
ExecStartPre=/bin/mkdir -p /root/.openclaw/workspace/projects/claimcoin-autoclaim/logs
ExecStart=/usr/bin/python3 -u -m claimcoin_autoclaim.cli run-loop --config accounts.yaml
Restart=always
RestartSec=5
StandardOutput=append:/root/.openclaw/workspace/projects/claimcoin-autoclaim/logs/run-loop-screen.log
StandardError=append:/root/.openclaw/workspace/projects/claimcoin-autoclaim/logs/run-loop-screen.log

[Install]
WantedBy=multi-user.target
