# ops/systemd/claimcoin-flaresolverr.service
[Unit]
Description=ClaimCoin FlareSolverr Helper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/state/flaresolverr-exp/src
Environment=HOME=/root
Environment=PORT=8195
Environment=HEADLESS=true
ExecStart=/root/.openclaw/workspace/state/flaresolverr-exp/src/.venv/bin/python src/flaresolverr.py
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
