# ops/systemd/claimcoin-antibot.service
[Unit]
Description=ClaimCoin AntiBot Solver API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/tmp-gh/antibot-image-solver
Environment=HOME=/root
Environment=ANTIBOT_OCR_PROFILE=fast
ExecStartPre=/bin/mkdir -p /root/.openclaw/workspace/tmp-gh/antibot-image-solver/logs
ExecStart=/root/.openclaw/workspace/tmp-gh/antibot-image-solver/.venv/bin/antibot-image-solver-api --host 127.0.0.1 --port 8010
Restart=always
RestartSec=2
StandardOutput=append:/root/.openclaw/workspace/tmp-gh/antibot-image-solver/logs/api.log
StandardError=append:/root/.openclaw/workspace/tmp-gh/antibot-image-solver/logs/api.log

[Install]
WantedBy=multi-user.target
