# ops/systemd/claimcoin-antibot.service
[Unit]
Description=ClaimCoin AntiBot Solver API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/projects/antibot-image-solver
Environment=HOME=/root
Environment=ANTIBOT_OCR_PROFILE=fast
ExecStart=/root/.openclaw/workspace/projects/antibot-image-solver/.venv/bin/antibot-image-solver-api --host 127.0.0.1 --port 8010
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
